# systemd/heating-brain.service
[Unit]
Description=Heating Brain (Tado controller)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=heating-brain
Group=heating-brain
WorkingDirectory=/opt/heating-brain
ExecStart=/opt/heating-brain/venv/bin/python -m app.orchestrator --config /etc/heating-brain/config.yaml
Restart=on-failure
RestartSec=10

# Security hardening
# NoNewPrivileges is intentionally NOT set here: the /api/update endpoint
# shells out via sudo to trigger the update oneshot service, which needs
# to gain privileges. The sudoers rule limits the service user to a
# single exact command so this remains tightly scoped.
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
# /opt/heating-brain must be writable so the /api/update endpoint can
# git fetch (writes FETCH_HEAD and objects). Without this, fetch fails
# with "Read-only file system" and the Update button silently returns
# up-to-date.
# /etc/heating-brain must be writable so the /api/schedule and /api/pin
# endpoints can rewrite config.yaml. Without this, saving a schedule
# window from the web UI fails with "Read-only file system".
ReadWritePaths=/var/lib/heating-brain /var/log/heating-brain /opt/heating-brain /etc/heating-brain

# Standard streams -> journal
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
